# power_supply@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/power_supply %I
