# powerdown.rules
ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/usr/bin/systemctl restart power_supply@true.service"
ENV{POWER_SUPPLY_ONLINE}=="1", RUN+="/usr/bin/systemctl restart power_supply@false.service"